# conf/prcdns-supervisor.conf
[program:prcdns]
command=prcdns --cn tcp/114.114.114/53,tcp/114.115.115/53
autostart=true
autorestart=true
user=focus
